# consumer.service
[Unit]
Description=SQS Consumer for ComfyUI video generation
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/venv/main/bin/python /workspace/consumer.py
WorkingDirectory=/workspace
EnvironmentFile=/workspace/consumer.env
Restart=always
RestartSec=5
StandardOutput=append:/workspace/consumer.log
StandardError=append:/workspace/consumer.log

[Install]
WantedBy=multi-user.target
