# consumer-supervisor.conf
[supervisord]
nodaemon=false
logfile=/workspace/supervisord.log
pidfile=/workspace/supervisord.pid

[program:consumer]
command=/venv/main/bin/python /workspace/consumer.py
directory=/workspace
autostart=true
autorestart=true
startsecs=5
stdout_logfile=/workspace/consumer.log
stderr_logfile=/workspace/consumer.log
